# logip/logip.service
# /etc/systemd/system/logip.service

[Unit]
Description=A script to log your public IP address.

[Service]
ExecStart=/usr/local/bin/logip.pl
